# docker/supervisord.conf
; Supervisor runs the three Coffer processes in one container:
;   - the FrankenPHP web server
;   - a single queue worker (SQLite is single-writer, so do not scale this up)
;   - the scheduler (prunes orphaned uploads, trashed files, expired downloads)
;
; If any program dies, the fatal event listener stops supervisor so the
; container exits and the orchestrator restarts it (fail-fast).

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:frankenphp]
command=frankenphp run --config /etc/frankenphp/Caddyfile
user=appuser
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue]
command=php /app/artisan queue:work --tries=3 --max-time=3600 --sleep=3
user=appuser
autostart=true
autorestart=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=php /app/artisan schedule:work
user=appuser
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:fail-fast]
command=bash -c "printf 'READY\n' && read line && kill -SIGQUIT $PPID"
events=PROCESS_STATE_FATAL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
